# build_files/rootfs/etc/xdg/autostart/luhood-welcome.desktop
[Desktop Entry]
Type=Application
Name=LuhoodOS Welcome
Exec=/usr/bin/luhood-welcome
OnlyShowIn=KDE;
X-KDE-AutostartScript=true
NoDisplay=true
